# hook/ld.so.preload
/home/red/open/github/tools/hook/libhook.so
